# packaging/systemd/gah-watchdog.service
[Unit]
Description=GAH Telegram health notification
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=khing
Group=khing
# systemd expands percent specifiers before invoking the shell. Escape the
# printf placeholder so `%s` reaches /bin/sh instead of becoming the user's
# login shell path.
ExecStart=/bin/sh -c 'msg=$$(/usr/bin/python3 /home/khing/workspace/agent-lab/scripts/gah-watchdog.py); if [ -n "$$msg" ]; then printf "%%s\\n" "$$msg" | /usr/local/bin/hermes send -t telegram -q; fi'
TimeoutStartSec=90
StandardOutput=journal
StandardError=journal
